# rootfs/var/lib/dpkg/info/ceph-common.list
/.
/etc
/etc/ceph
/etc/ceph/rbdmap
/etc/default
/etc/default/ceph
/etc/init.d
/etc/init.d/rbdmap
/etc/logrotate.d
/etc/logrotate.d/ceph-common
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ceph.target
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/50-rbd.rules
/usr
/usr/bin
/usr/bin/ceph
/usr/bin/ceph-authtool
/usr/bin/ceph-brag
/usr/bin/ceph-conf
/usr/bin/ceph-crush-location
/usr/bin/ceph-dencoder
/usr/bin/ceph-post-file
/usr/bin/ceph-rbdnamer
/usr/bin/ceph-syn
/usr/bin/rados
/usr/bin/rbd
/usr/bin/rbd-replay
/usr/bin/rbd-replay-many
/usr/bin/rbd-replay-prep
/usr/bin/rbdmap
/usr/lib
/usr/lib/python2.7
/usr/lib/python2.7/dist-packages
/usr/lib/python2.7/dist-packages/ceph_argparse.py
/usr/lib/python2.7/dist-packages/ceph_daemon.py
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/ceph.conf
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/rados
/usr/share/bash-completion/completions/rbd
/usr/share/ceph
/usr/share/ceph/id_rsa_drop.ceph.com
/usr/share/ceph/id_rsa_drop.ceph.com.pub
/usr/share/ceph/known_hosts_drop.ceph.com
/usr/share/doc
/usr/share/doc/ceph-common
/usr/share/doc/ceph-common/changelog.Debian.gz
/usr/share/doc/ceph-common/changelog.gz
/usr/share/doc/ceph-common/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/ceph-common
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/ceph-crush-location.1.gz
/usr/share/man/man8
/usr/share/man/man8/ceph-authtool.8.gz
/usr/share/man/man8/ceph-conf.8.gz
/usr/share/man/man8/ceph-dencoder.8.gz
/usr/share/man/man8/ceph-post-file.8.gz
/usr/share/man/man8/ceph-rbdnamer.8.gz
/usr/share/man/man8/ceph-syn.8.gz
/usr/share/man/man8/ceph.8.gz
/usr/share/man/man8/rados.8.gz
/usr/share/man/man8/rbd-replay-many.8.gz
/usr/share/man/man8/rbd-replay-prep.8.gz
/usr/share/man/man8/rbd-replay.8.gz
/usr/share/man/man8/rbd.8.gz
/usr/share/man/man8/rbdmap.8.gz
/var
/var/lib
/var/lib/ceph
/var/log
/var/log/ceph
