# rootfs/var/lib/dpkg/info/ksm-control-daemon.list
/.
/etc
/etc/default
/etc/default/ksmtuned
/etc/ksmtuned.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ksmtuned.service
/usr
/usr/sbin
/usr/sbin/ksmtuned
/usr/share
/usr/share/doc
/usr/share/doc/ksm-control-daemon
/usr/share/doc/ksm-control-daemon/SOURCE
/usr/share/doc/ksm-control-daemon/changelog.Debian.gz
/usr/share/doc/ksm-control-daemon/copyright
