# rootfs/var/lib/dpkg/info/lxc-pve.list
/.
/etc
/etc/apparmor.d
/etc/apparmor.d/abstractions
/etc/apparmor.d/abstractions/lxc
/etc/apparmor.d/abstractions/lxc/container-base
/etc/apparmor.d/abstractions/lxc/start-container
/etc/apparmor.d/lxc
/etc/apparmor.d/lxc/lxc-default
/etc/apparmor.d/lxc/lxc-default-cgns
/etc/apparmor.d/lxc/lxc-default-with-mounting
/etc/apparmor.d/lxc/lxc-default-with-nesting
/etc/apparmor.d/lxc-containers
/etc/apparmor.d/usr.bin.lxc-start
/etc/default
/etc/default/lxc
/etc/lxc
/etc/lxc/default.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/lxc-monitord.service
/lib/systemd/system/lxc-net.service
/lib/systemd/system/lxc.service
/lib/systemd/system/lxc@.service
/usr
/usr/bin
/usr/bin/lxc-attach
/usr/bin/lxc-autostart
/usr/bin/lxc-cgroup
/usr/bin/lxc-checkconfig
/usr/bin/lxc-checkpoint
/usr/bin/lxc-config
/usr/bin/lxc-console
/usr/bin/lxc-copy
/usr/bin/lxc-create
/usr/bin/lxc-destroy
/usr/bin/lxc-device
/usr/bin/lxc-execute
/usr/bin/lxc-freeze
/usr/bin/lxc-info
/usr/bin/lxc-ls
/usr/bin/lxc-monitor
/usr/bin/lxc-snapshot
/usr/bin/lxc-start
/usr/bin/lxc-stop
/usr/bin/lxc-top
/usr/bin/lxc-unfreeze
/usr/bin/lxc-unshare
/usr/bin/lxc-update-config
/usr/bin/lxc-usernsexec
/usr/bin/lxc-wait
/usr/lib
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/liblxc.so.1.4.0
/usr/lib/x86_64-linux-gnu/lxc
/usr/lib/x86_64-linux-gnu/lxc/hooks
/usr/lib/x86_64-linux-gnu/lxc/hooks/unmount-namespace
/usr/lib/x86_64-linux-gnu/lxc/lxc-apparmor-load
/usr/lib/x86_64-linux-gnu/lxc/lxc-containers
/usr/lib/x86_64-linux-gnu/lxc/lxc-monitord
/usr/lib/x86_64-linux-gnu/lxc/lxc-net
/usr/lib/x86_64-linux-gnu/lxc/lxc-user-nic
/usr/lib/x86_64-linux-gnu/lxc/rootfs
/usr/lib/x86_64-linux-gnu/lxc/rootfs/README
/usr/sbin
/usr/sbin/init.lxc
/usr/sbin/init.lxc.static
/usr/share
/usr/share/apparmor-features
/usr/share/apparmor-features/features
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/lxc
/usr/share/doc
/usr/share/doc/lxc-pve
/usr/share/doc/lxc-pve/changelog.Debian.gz
/usr/share/doc/lxc-pve/copyright
/usr/share/lxc
/usr/share/lxc/config
/usr/share/lxc/config/alpine.common.conf
/usr/share/lxc/config/alpine.userns.conf
/usr/share/lxc/config/archlinux.common.conf
/usr/share/lxc/config/archlinux.userns.conf
/usr/share/lxc/config/centos.common.conf
/usr/share/lxc/config/centos.userns.conf
/usr/share/lxc/config/common.conf
/usr/share/lxc/config/common.conf.d
/usr/share/lxc/config/common.conf.d/README
/usr/share/lxc/config/common.seccomp
/usr/share/lxc/config/debian.common.conf
/usr/share/lxc/config/debian.userns.conf
/usr/share/lxc/config/fedora.common.conf
/usr/share/lxc/config/fedora.userns.conf
/usr/share/lxc/config/gentoo.common.conf
/usr/share/lxc/config/gentoo.moresecure.conf
/usr/share/lxc/config/gentoo.userns.conf
/usr/share/lxc/config/nesting.conf
/usr/share/lxc/config/oci.common.conf
/usr/share/lxc/config/opensuse.common.conf
/usr/share/lxc/config/opensuse.userns.conf
/usr/share/lxc/config/openwrt.common.conf
/usr/share/lxc/config/oracle.common.conf
/usr/share/lxc/config/oracle.userns.conf
/usr/share/lxc/config/plamo.common.conf
/usr/share/lxc/config/plamo.userns.conf
/usr/share/lxc/config/sabayon.common.conf
/usr/share/lxc/config/sabayon.userns.conf
/usr/share/lxc/config/slackware.common.conf
/usr/share/lxc/config/slackware.userns.conf
/usr/share/lxc/config/sparclinux.common.conf
/usr/share/lxc/config/sparclinux.userns.conf
/usr/share/lxc/config/ubuntu-cloud.common.conf
/usr/share/lxc/config/ubuntu-cloud.lucid.conf
/usr/share/lxc/config/ubuntu-cloud.userns.conf
/usr/share/lxc/config/ubuntu.common.conf
/usr/share/lxc/config/ubuntu.lucid.conf
/usr/share/lxc/config/ubuntu.userns.conf
/usr/share/lxc/config/userns.conf
/usr/share/lxc/config/voidlinux.common.conf
/usr/share/lxc/config/voidlinux.userns.conf
/usr/share/lxc/hooks
/usr/share/lxc/hooks/clonehostname
/usr/share/lxc/hooks/dhclient
/usr/share/lxc/hooks/dhclient-script
/usr/share/lxc/hooks/mountecryptfsroot
/usr/share/lxc/hooks/nvidia
/usr/share/lxc/hooks/squid-deb-proxy-client
/usr/share/lxc/hooks/ubuntu-cloud-prep
/usr/share/lxc/lxc-patch.py
/usr/share/lxc/lxc.functions
/usr/share/lxc/selinux
/usr/share/lxc/selinux/lxc.if
/usr/share/lxc/selinux/lxc.te
/usr/share/lxc/templates
/usr/share/lxc/templates/lxc-busybox
/usr/share/lxc/templates/lxc-download
/usr/share/lxc/templates/lxc-local
/usr/share/lxc/templates/lxc-oci
/usr/share/man
/usr/share/man/ja
/usr/share/man/ja/man1
/usr/share/man/ja/man1/lxc-attach.1.gz
/usr/share/man/ja/man1/lxc-autostart.1.gz
/usr/share/man/ja/man1/lxc-cgroup.1.gz
/usr/share/man/ja/man1/lxc-checkconfig.1.gz
/usr/share/man/ja/man1/lxc-checkpoint.1.gz
/usr/share/man/ja/man1/lxc-config.1.gz
/usr/share/man/ja/man1/lxc-console.1.gz
/usr/share/man/ja/man1/lxc-copy.1.gz
/usr/share/man/ja/man1/lxc-create.1.gz
/usr/share/man/ja/man1/lxc-destroy.1.gz
/usr/share/man/ja/man1/lxc-device.1.gz
/usr/share/man/ja/man1/lxc-execute.1.gz
/usr/share/man/ja/man1/lxc-freeze.1.gz
/usr/share/man/ja/man1/lxc-info.1.gz
/usr/share/man/ja/man1/lxc-ls.1.gz
/usr/share/man/ja/man1/lxc-monitor.1.gz
/usr/share/man/ja/man1/lxc-snapshot.1.gz
/usr/share/man/ja/man1/lxc-start.1.gz
/usr/share/man/ja/man1/lxc-stop.1.gz
/usr/share/man/ja/man1/lxc-top.1.gz
/usr/share/man/ja/man1/lxc-unfreeze.1.gz
/usr/share/man/ja/man1/lxc-unshare.1.gz
/usr/share/man/ja/man1/lxc-update-config.1.gz
/usr/share/man/ja/man1/lxc-user-nic.1.gz
/usr/share/man/ja/man1/lxc-usernsexec.1.gz
/usr/share/man/ja/man1/lxc-wait.1.gz
/usr/share/man/ja/man5
/usr/share/man/ja/man5/lxc-usernet.5.gz
/usr/share/man/ja/man5/lxc.conf.5.gz
/usr/share/man/ja/man5/lxc.container.conf.5.gz
/usr/share/man/ja/man5/lxc.system.conf.5.gz
/usr/share/man/ja/man7
/usr/share/man/ja/man7/lxc.7.gz
/usr/share/man/ko
/usr/share/man/ko/man1
/usr/share/man/ko/man1/lxc-attach.1.gz
/usr/share/man/ko/man1/lxc-autostart.1.gz
/usr/share/man/ko/man1/lxc-cgroup.1.gz
/usr/share/man/ko/man1/lxc-checkconfig.1.gz
/usr/share/man/ko/man1/lxc-checkpoint.1.gz
/usr/share/man/ko/man1/lxc-config.1.gz
/usr/share/man/ko/man1/lxc-console.1.gz
/usr/share/man/ko/man1/lxc-copy.1.gz
/usr/share/man/ko/man1/lxc-create.1.gz
/usr/share/man/ko/man1/lxc-destroy.1.gz
/usr/share/man/ko/man1/lxc-device.1.gz
/usr/share/man/ko/man1/lxc-execute.1.gz
/usr/share/man/ko/man1/lxc-freeze.1.gz
/usr/share/man/ko/man1/lxc-info.1.gz
/usr/share/man/ko/man1/lxc-ls.1.gz
/usr/share/man/ko/man1/lxc-monitor.1.gz
/usr/share/man/ko/man1/lxc-snapshot.1.gz
/usr/share/man/ko/man1/lxc-start.1.gz
/usr/share/man/ko/man1/lxc-stop.1.gz
/usr/share/man/ko/man1/lxc-top.1.gz
/usr/share/man/ko/man1/lxc-unfreeze.1.gz
/usr/share/man/ko/man1/lxc-unshare.1.gz
/usr/share/man/ko/man1/lxc-user-nic.1.gz
/usr/share/man/ko/man1/lxc-usernsexec.1.gz
/usr/share/man/ko/man1/lxc-wait.1.gz
/usr/share/man/ko/man5
/usr/share/man/ko/man5/lxc-usernet.5.gz
/usr/share/man/ko/man5/lxc.conf.5.gz
/usr/share/man/ko/man5/lxc.container.conf.5.gz
/usr/share/man/ko/man5/lxc.system.conf.5.gz
/usr/share/man/ko/man7
/usr/share/man/ko/man7/lxc.7.gz
/usr/share/man/man1
/usr/share/man/man1/lxc-attach.1.gz
/usr/share/man/man1/lxc-autostart.1.gz
/usr/share/man/man1/lxc-cgroup.1.gz
/usr/share/man/man1/lxc-checkconfig.1.gz
/usr/share/man/man1/lxc-checkpoint.1.gz
/usr/share/man/man1/lxc-config.1.gz
/usr/share/man/man1/lxc-console.1.gz
/usr/share/man/man1/lxc-copy.1.gz
/usr/share/man/man1/lxc-create.1.gz
/usr/share/man/man1/lxc-destroy.1.gz
/usr/share/man/man1/lxc-device.1.gz
/usr/share/man/man1/lxc-execute.1.gz
/usr/share/man/man1/lxc-freeze.1.gz
/usr/share/man/man1/lxc-info.1.gz
/usr/share/man/man1/lxc-ls.1.gz
/usr/share/man/man1/lxc-monitor.1.gz
/usr/share/man/man1/lxc-snapshot.1.gz
/usr/share/man/man1/lxc-start.1.gz
/usr/share/man/man1/lxc-stop.1.gz
/usr/share/man/man1/lxc-top.1.gz
/usr/share/man/man1/lxc-unfreeze.1.gz
/usr/share/man/man1/lxc-unshare.1.gz
/usr/share/man/man1/lxc-update-config.1.gz
/usr/share/man/man1/lxc-user-nic.1.gz
/usr/share/man/man1/lxc-usernsexec.1.gz
/usr/share/man/man1/lxc-wait.1.gz
/usr/share/man/man5
/usr/share/man/man5/lxc-usernet.5.gz
/usr/share/man/man5/lxc.conf.5.gz
/usr/share/man/man5/lxc.container.conf.5.gz
/usr/share/man/man5/lxc.system.conf.5.gz
/usr/share/man/man7
/usr/share/man/man7/lxc.7.gz
/var
/var/cache
/var/cache/lxc
/var/lib
/var/lib/lxc
/var/log
/var/log/lxc
/usr/lib/x86_64-linux-gnu/liblxc.so.1
