# rootfs/var/lib/dpkg/info/pve-container.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/lxc@.service.d
/lib/systemd/system/lxc@.service.d/pve-reboot.conf
/lib/systemd/system/pve-container@.service
/lib/systemd/system/system-pve\x2dcontainer.slice
/usr
/usr/sbin
/usr/sbin/pct
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/pct
/usr/share/doc
/usr/share/doc/pve-container
/usr/share/doc/pve-container/changelog.Debian.gz
/usr/share/doc/pve-container/copyright
/usr/share/lxc
/usr/share/lxc/config
/usr/share/lxc/config/common.conf.d
/usr/share/lxc/config/common.conf.d/01-pve.conf
/usr/share/lxc/config/pve-userns.seccomp
/usr/share/lxc/hooks
/usr/share/lxc/hooks/lxc-pve-autodev-hook
/usr/share/lxc/hooks/lxc-pve-poststop-hook
/usr/share/lxc/hooks/lxc-pve-prestart-hook
/usr/share/lxc/lxc-pve-reboot-trigger
/usr/share/lxc/lxcnetaddbr
/usr/share/lxc/pve-container-stop-wrapper
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/pct.1.gz
/usr/share/man/man5
/usr/share/man/man5/pct.conf.5.gz
/usr/share/perl5
/usr/share/perl5/PVE
/usr/share/perl5/PVE/API2
/usr/share/perl5/PVE/API2/LXC
/usr/share/perl5/PVE/API2/LXC/Config.pm
/usr/share/perl5/PVE/API2/LXC/Snapshot.pm
/usr/share/perl5/PVE/API2/LXC/Status.pm
/usr/share/perl5/PVE/API2/LXC.pm
/usr/share/perl5/PVE/CLI
/usr/share/perl5/PVE/CLI/pct.pm
/usr/share/perl5/PVE/LXC
/usr/share/perl5/PVE/LXC/Config.pm
/usr/share/perl5/PVE/LXC/Create.pm
/usr/share/perl5/PVE/LXC/Migrate.pm
/usr/share/perl5/PVE/LXC/Setup
/usr/share/perl5/PVE/LXC/Setup/Alpine.pm
/usr/share/perl5/PVE/LXC/Setup/ArchLinux.pm
/usr/share/perl5/PVE/LXC/Setup/Base.pm
/usr/share/perl5/PVE/LXC/Setup/CentOS.pm
/usr/share/perl5/PVE/LXC/Setup/Debian.pm
/usr/share/perl5/PVE/LXC/Setup/Devuan.pm
/usr/share/perl5/PVE/LXC/Setup/Fedora.pm
/usr/share/perl5/PVE/LXC/Setup/Gentoo.pm
/usr/share/perl5/PVE/LXC/Setup/SUSE.pm
/usr/share/perl5/PVE/LXC/Setup/Ubuntu.pm
/usr/share/perl5/PVE/LXC/Setup.pm
/usr/share/perl5/PVE/LXC.pm
/usr/share/perl5/PVE/VZDump
/usr/share/perl5/PVE/VZDump/ConvertOVZ.pm
/usr/share/perl5/PVE/VZDump/LXC.pm
/usr/share/man/man5/ct.conf.5.gz
