# rootfs/var/lib/dpkg/info/pve-cluster.list
/.
/etc
/etc/sysctl.d
/etc/sysctl.d/pve.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/pve-cluster.service
/usr
/usr/bin
/usr/bin/create_pmxcfs_db
/usr/bin/pmxcfs
/usr/bin/pvecm
/usr/lib
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/perl5
/usr/lib/x86_64-linux-gnu/perl5/5.24
/usr/lib/x86_64-linux-gnu/perl5/5.24/auto
/usr/lib/x86_64-linux-gnu/perl5/5.24/auto/PVE
/usr/lib/x86_64-linux-gnu/perl5/5.24/auto/PVE/IPCC
/usr/lib/x86_64-linux-gnu/perl5/5.24/auto/PVE/IPCC/IPCC.so
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/pvecm
/usr/share/doc
/usr/share/doc/pve-cluster
/usr/share/doc/pve-cluster/README.gz
/usr/share/doc/pve-cluster/SOURCE
/usr/share/doc/pve-cluster/changelog.Debian.gz
/usr/share/doc/pve-cluster/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/pvecm.1.gz
/usr/share/man/man5
/usr/share/man/man5/datacenter.cfg.5.gz
/usr/share/man/man8
/usr/share/man/man8/pmxcfs.8.gz
/usr/share/perl5
/usr/share/perl5/PVE
/usr/share/perl5/PVE/API2
/usr/share/perl5/PVE/API2/ClusterConfig.pm
/usr/share/perl5/PVE/CLI
/usr/share/perl5/PVE/CLI/pvecm.pm
/usr/share/perl5/PVE/Cluster.pm
/usr/share/perl5/PVE/Corosync.pm
/usr/share/perl5/PVE/IPCC.pm
/var
/var/lib
/var/lib/pve-cluster
