# rootfs/var/lib/dpkg/info/pve-firewall.list
/.
/etc
/etc/logrotate.d
/etc/logrotate.d/pve-firewall
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/pve-firewall.service
/lib/systemd/system/pvefw-logger.service
/usr
/usr/lib
/usr/lib/sysctl.d
/usr/lib/sysctl.d/pve-firewall.conf
/usr/sbin
/usr/sbin/pve-firewall
/usr/sbin/pvefw-logger
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/pve-firewall
/usr/share/doc
/usr/share/doc/pve-firewall
/usr/share/doc/pve-firewall/README
/usr/share/doc/pve-firewall/SOURCE
/usr/share/doc/pve-firewall/changelog.Debian.gz
/usr/share/doc/pve-firewall/copyright
/usr/share/doc/pve-firewall/example
/usr/share/doc/pve-firewall/example/100.fw
/usr/share/doc/pve-firewall/example/cluster.fw
/usr/share/doc/pve-firewall/example/host.fw
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/pve-firewall.8.gz
/usr/share/perl5
/usr/share/perl5/PVE
/usr/share/perl5/PVE/API2
/usr/share/perl5/PVE/API2/Firewall
/usr/share/perl5/PVE/API2/Firewall/Aliases.pm
/usr/share/perl5/PVE/API2/Firewall/Cluster.pm
/usr/share/perl5/PVE/API2/Firewall/Groups.pm
/usr/share/perl5/PVE/API2/Firewall/Host.pm
/usr/share/perl5/PVE/API2/Firewall/IPSet.pm
/usr/share/perl5/PVE/API2/Firewall/Rules.pm
/usr/share/perl5/PVE/API2/Firewall/VM.pm
/usr/share/perl5/PVE/Firewall.pm
/usr/share/perl5/PVE/FirewallSimulator.pm
/usr/share/perl5/PVE/Service
/usr/share/perl5/PVE/Service/pve_firewall.pm
/var
/var/lib
/var/lib/pve-firewall
