# rootfs/var/lib/dpkg/info/rrdcached.list
/.
/etc
/etc/default
/etc/default/rrdcached
/etc/init.d
/etc/init.d/rrdcached
/usr
/usr/bin
/usr/bin/rrdcached
/usr/share
/usr/share/doc
/usr/share/doc/rrdcached
/usr/share/doc/rrdcached/CONTRIBUTORS.gz
/usr/share/doc/rrdcached/NEWS.Debian.gz
/usr/share/doc/rrdcached/NEWS.gz
/usr/share/doc/rrdcached/changelog.Debian.amd64.gz
/usr/share/doc/rrdcached/changelog.Debian.gz
/usr/share/doc/rrdcached/changelog.gz
/usr/share/doc/rrdcached/copyright
/usr/share/doc/rrdcached/examples
/usr/share/doc/rrdcached/examples/RRDCached.pm
/usr/share/doc/rrdcached/examples/rrdcached-size.pl
/usr/share/doc/rrdcached/html
/usr/share/doc/rrdcached/html/rrdcached.html
/usr/share/doc-base
/usr/share/doc-base/rrdcached
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/rrdcached.1.gz
