# rootfs/var/lib/dpkg/info/rrdcached.conffiles
/etc/default/rrdcached
/etc/init.d/rrdcached
