# rootfs/var/lib/dpkg/info/zfsutils-linux.list
/.
/etc
/etc/cron.d
/etc/cron.d/zfsutils-linux
/etc/default
/etc/default/zfs
/etc/sudoers.d
/etc/sudoers.d/zfs
/etc/zfs
/etc/zfs/zfs-functions
/etc/zfs/zpool.d
/lib
/lib/modules-load.d
/lib/modules-load.d/zfs.conf
/lib/systemd
/lib/systemd/system
/lib/systemd/system/zfs-import-cache.service
/lib/systemd/system/zfs-import-scan.service
/lib/systemd/system/zfs-import.target
/lib/systemd/system/zfs-mount.service
/lib/systemd/system/zfs-share.service
/lib/systemd/system/zfs.target
/lib/systemd/system-preset
/lib/systemd/system-preset/50-zfs.preset
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/60-zvol.rules
/lib/udev/rules.d/69-vdev.rules
/lib/udev/rules.d/90-zfs.rules
/lib/udev/vdev_id
/lib/udev/zvol_id
/sbin
/sbin/fsck.zfs
/sbin/mount.zfs
/sbin/zdb
/sbin/zfs
/sbin/zhack
/sbin/zinject
/sbin/zpios
/sbin/zpool
/sbin/zstreamdump
/sbin/ztest
/usr
/usr/bin
/usr/bin/zgenhostid
/usr/lib
/usr/lib/os-probes
/usr/lib/os-probes/10zvol-test
/usr/lib/zfs-linux
/usr/lib/zfs-linux/scrub
/usr/lib/zfs-linux/zpool.d
/usr/lib/zfs-linux/zpool.d/ata_err
/usr/lib/zfs-linux/zpool.d/cmd_to
/usr/lib/zfs-linux/zpool.d/defect
/usr/lib/zfs-linux/zpool.d/enc
/usr/lib/zfs-linux/zpool.d/encdev
/usr/lib/zfs-linux/zpool.d/fault_led
/usr/lib/zfs-linux/zpool.d/health
/usr/lib/zfs-linux/zpool.d/hours_on
/usr/lib/zfs-linux/zpool.d/iostat
/usr/lib/zfs-linux/zpool.d/iostat-10s
/usr/lib/zfs-linux/zpool.d/iostat-1s
/usr/lib/zfs-linux/zpool.d/label
/usr/lib/zfs-linux/zpool.d/locate_led
/usr/lib/zfs-linux/zpool.d/lsblk
/usr/lib/zfs-linux/zpool.d/media
/usr/lib/zfs-linux/zpool.d/model
/usr/lib/zfs-linux/zpool.d/nonmed
/usr/lib/zfs-linux/zpool.d/nvme_err
/usr/lib/zfs-linux/zpool.d/off_ucor
/usr/lib/zfs-linux/zpool.d/pend_sec
/usr/lib/zfs-linux/zpool.d/pwr_cyc
/usr/lib/zfs-linux/zpool.d/r_proc
/usr/lib/zfs-linux/zpool.d/r_ucor
/usr/lib/zfs-linux/zpool.d/realloc
/usr/lib/zfs-linux/zpool.d/rep_ucor
/usr/lib/zfs-linux/zpool.d/serial
/usr/lib/zfs-linux/zpool.d/ses
/usr/lib/zfs-linux/zpool.d/size
/usr/lib/zfs-linux/zpool.d/slaves
/usr/lib/zfs-linux/zpool.d/slot
/usr/lib/zfs-linux/zpool.d/smart
/usr/lib/zfs-linux/zpool.d/smart_test
/usr/lib/zfs-linux/zpool.d/smartx
/usr/lib/zfs-linux/zpool.d/temp
/usr/lib/zfs-linux/zpool.d/test_ended
/usr/lib/zfs-linux/zpool.d/test_progress
/usr/lib/zfs-linux/zpool.d/test_status
/usr/lib/zfs-linux/zpool.d/test_type
/usr/lib/zfs-linux/zpool.d/upath
/usr/lib/zfs-linux/zpool.d/vendor
/usr/lib/zfs-linux/zpool.d/w_proc
/usr/lib/zfs-linux/zpool.d/w_ucor
/usr/sbin
/usr/sbin/arc_summary
/usr/sbin/arcstat
/usr/sbin/dbufstat
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/zfs
/usr/share/doc
/usr/share/doc/zfsutils-linux
/usr/share/doc/zfsutils-linux/COPYRIGHT
/usr/share/doc/zfsutils-linux/OPENSOLARIS.LICENSE.gz
/usr/share/doc/zfsutils-linux/changelog.Debian.gz
/usr/share/doc/zfsutils-linux/copyright
/usr/share/doc/zfsutils-linux/examples
/usr/share/doc/zfsutils-linux/examples/vdev_id.conf.alias.example
/usr/share/doc/zfsutils-linux/examples/vdev_id.conf.multipath.example
/usr/share/doc/zfsutils-linux/examples/vdev_id.conf.sas_direct.example
/usr/share/doc/zfsutils-linux/examples/vdev_id.conf.sas_switch.example
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/zdev
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/raidz_test.1.gz
/usr/share/man/man1/test-runner.1.gz
/usr/share/man/man1/zhack.1.gz
/usr/share/man/man1/zpios.1.gz
/usr/share/man/man1/ztest.1.gz
/usr/share/man/man5
/usr/share/man/man5/vdev_id.conf.5.gz
/usr/share/man/man5/zfs-events.5.gz
/usr/share/man/man5/zfs-module-parameters.5.gz
/usr/share/man/man5/zpool-features.5.gz
/usr/share/man/man8
/usr/share/man/man8/fsck.zfs.8.gz
/usr/share/man/man8/mount.zfs.8.gz
/usr/share/man/man8/vdev_id.8.gz
/usr/share/man/man8/zdb.8.gz
/usr/share/man/man8/zfs.8.gz
/usr/share/man/man8/zgenhostid.8.gz
/usr/share/man/man8/zinject.8.gz
/usr/share/man/man8/zpool.8.gz
/usr/share/man/man8/zstreamdump.8.gz
/etc/zfs/zpool.d/ata_err
/etc/zfs/zpool.d/cmd_to
/etc/zfs/zpool.d/defect
/etc/zfs/zpool.d/enc
/etc/zfs/zpool.d/encdev
/etc/zfs/zpool.d/fault_led
/etc/zfs/zpool.d/health
/etc/zfs/zpool.d/hours_on
/etc/zfs/zpool.d/iostat
/etc/zfs/zpool.d/iostat-10s
/etc/zfs/zpool.d/iostat-1s
/etc/zfs/zpool.d/label
/etc/zfs/zpool.d/locate_led
/etc/zfs/zpool.d/lsblk
/etc/zfs/zpool.d/media
/etc/zfs/zpool.d/model
/etc/zfs/zpool.d/nonmed
/etc/zfs/zpool.d/nvme_err
/etc/zfs/zpool.d/off_ucor
/etc/zfs/zpool.d/pend_sec
/etc/zfs/zpool.d/pwr_cyc
/etc/zfs/zpool.d/r_proc
/etc/zfs/zpool.d/r_ucor
/etc/zfs/zpool.d/realloc
/etc/zfs/zpool.d/rep_ucor
/etc/zfs/zpool.d/serial
/etc/zfs/zpool.d/ses
/etc/zfs/zpool.d/size
/etc/zfs/zpool.d/slaves
/etc/zfs/zpool.d/slot
/etc/zfs/zpool.d/smart
/etc/zfs/zpool.d/smart_test
/etc/zfs/zpool.d/smartx
/etc/zfs/zpool.d/temp
/etc/zfs/zpool.d/test_ended
/etc/zfs/zpool.d/test_progress
/etc/zfs/zpool.d/test_status
/etc/zfs/zpool.d/test_type
/etc/zfs/zpool.d/upath
/etc/zfs/zpool.d/vendor
/etc/zfs/zpool.d/w_proc
/etc/zfs/zpool.d/w_ucor
