# rootfs/var/lib/dpkg/info/zfsutils-linux.conffiles
/etc/cron.d/zfsutils-linux
/etc/default/zfs
/etc/sudoers.d/zfs
/etc/zfs/zfs-functions
